# automation/daemon/com.ghostmonday.director.stage-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghostmonday.director.stage-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/user944529/Desktop/director/automation/daemon/stage-validation-daemon.sh</string>
        <string>--continuous</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/user944529/Desktop/director/automation/logs/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user944529/Desktop/director/automation/logs/daemon-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/user944529/Desktop/director</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>AbandonProcessGroup</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
